# templates/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:apache2]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/apache2ctl -DFOREGROUND 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [APACHE] $line"; done'
user=root
autostart=true
autorestart=true
priority=100
stdout_logfile=/data/logs/apache2.log
stderr_logfile=/data/logs/apache2_error.log

[program:tftpd-hpa]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/in.tftpd --foreground --secure --create --address 0.0.0.0:69 /tftpboot 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [TFTP] $line"; done'
user=root
autostart=true
autorestart=true
priority=200
stdout_logfile=/data/logs/tftpd.log
stderr_logfile=/data/logs/tftpd_error.log

[program:rpcbind]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/rpcbind -f 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [RPCBIND] $line"; done'
user=root
autostart=true
autorestart=true
priority=300
stdout_logfile=/data/logs/rpcbind.log
stderr_logfile=/data/logs/rpcbind_error.log

[program:rpc-statd]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/rpc.statd --no-notify 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [RPC-STATD] $line"; done'
user=root
autostart=true
autorestart=false
exitcodes=0
priority=310
stdout_logfile=/data/logs/rpc-statd.log
stderr_logfile=/data/logs/rpc-statd_error.log

[program:nfs-kernel-server]
command=/bin/bash -c 'sleep 3 && TZ=$TZ exportfs -ra && /usr/sbin/rpc.nfsd 8 && /usr/sbin/rpc.mountd -F 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [NFS] $line"; done'
user=root
autostart=true
autorestart=true
priority=320
startsecs=5
stdout_logfile=/data/logs/nfs.log
stderr_logfile=/data/logs/nfs_error.log

[program:vsftpd]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/vsftpd /etc/vsftpd.conf 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [FTP] $line"; done'
user=root
autostart=true
autorestart=true
priority=325
stdout_logfile=/data/logs/vsftpd.log
stderr_logfile=/data/logs/vsftpd_error.log

[program:fog-image-replicator]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGImageReplicator/FOGImageReplicator 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [IMAGE-REPL] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/fogreplicator.log
stderr_logfile=/data/logs/fogreplicator.log

[program:fog-snapin-replicator]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGSnapinReplicator/FOGSnapinReplicator 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [SNAPIN-REPL] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/fogsnapinrep.log
stderr_logfile=/data/logs/fogsnapinrep.log

[program:fog-multicast-manager]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGMulticastManager/FOGMulticastManager 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [MULTICAST] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/multicast.log
stderr_logfile=/data/logs/multicast.log

[program:fog-task-scheduler]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGTaskScheduler/FOGTaskScheduler 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [TASK-SCHED] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/fogscheduler.log
stderr_logfile=/data/logs/fogscheduler.log

[program:fog-ping-hosts]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGPingHosts/FOGPingHosts 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [PING-HOSTS] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/pinghost.log
stderr_logfile=/data/logs/pinghost.log

[program:fog-snapin-hash]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGSnapinHash/FOGSnapinHash 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [SNAPIN-HASH] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/fogsnapinhash.log
stderr_logfile=/data/logs/fogsnapinhash.log

[program:fog-image-size]
command=/bin/bash -c 'TZ=$TZ php /opt/fog/service/FOGImageSize/FOGImageSize 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [IMAGE-SIZE] $line"; done'
user=www-data
autostart=true
autorestart=true
priority=400
stdout_logfile=/data/logs/fogimagesize.log
stderr_logfile=/data/logs/fogimagesize.log

[program:isc-dhcp-server]
command=/bin/bash -c 'TZ=$TZ /usr/sbin/dhcpd -f -d -cf /etc/dhcp/dhcpd.conf 2>&1 | while IFS= read -r line; do echo "[$(TZ=$TZ date "+%%Y-%%m-%%d %%H:%%M:%%S")] [DHCP] $line"; done'
user=root
autostart=true
autorestart=true
priority=350
stdout_logfile=/data/logs/dhcpd.log
stderr_logfile=/data/logs/dhcpd_error.log